# crontab
0 * * * * rclone sync "Mega Encrypt:Main" /home/doodles/Documents/Main && rclon>

